# openvpn-ui/openvpn-ui.service
[Unit]
Description=OpenVPN Management UI
Documentation=https://github.com/yourusername/openvpn-ui
After=network.target

[Service]
Type=simple
User=USERNAME
Group=USERNAME
WorkingDirectory=/home/USERNAME/openvpn-ui

# Path to Node.js
ExecStart=/usr/bin/node server.js

# Restart configuration
Restart=always
RestartSec=10

# Environment variables
Environment=NODE_ENV=production
Environment=PORT=8080
Environment=HOST=0.0.0.0
Environment=ADMIN_USERNAME=admin
Environment=ADMIN_PASSWORD=hprogramist8060
Environment=SESSION_SECRET=super-secret-change-me-in-production
Environment=OVPN_CONFIG_DIR=/home/USERNAME/ovpns
Environment=SYSTEM_USER=USERNAME
Environment=OPENVPN_SERVICE=openvpn

# Security settings
NoNewPrivileges=true
PrivateTmp=true

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=openvpn-ui

[Install]
WantedBy=multi-user.target
